# streamdeck-ctrl.service.in
# streamdeck-ctrl.service — generated by setup.sh from this template
#
# This service is NOT started at boot by default. It is started by the
# udev rule (99-streamdeck.rules) only when an Elgato Stream Deck is
# present. To start at boot regardless, run:
#   systemctl enable streamdeck-ctrl
#
[Unit]
Description=Stream Deck Control Daemon
After=multi-user.target

[Service]
Type=simple
User={USER}
Group=plugdev
WorkingDirectory={INSTALL_DIR}
Environment=PYTHONPATH={INSTALL_DIR}
ExecStart=/usr/bin/python3 {INSTALL_DIR}/streamdeck_ctrl/main.py --config={CONFIG_PATH}
ExecStartPost=-{INSTALL_DIR}/screens/display-control/scripts/sync-als-dimmer.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
RuntimeDirectory=streamdeck-ctrl
RuntimeDirectoryMode=0755

# No [Install] section — service is triggered by udev, not enabled at boot
