# streamdeck-ctrl.service.in
# streamdeck-ctrl.service — generated by setup.sh from this template
[Unit]
Description=Stream Deck Control Daemon
After=multi-user.target
# udev TAG+="systemd" also triggers this unit on USB plug-in

[Service]
Type=simple
User={USER}
Group=plugdev
WorkingDirectory={INSTALL_DIR}
ExecStart=/usr/bin/python3 {INSTALL_DIR}/streamdeck_ctrl/main.py --config={CONFIG_PATH}
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
RuntimeDirectory=streamdeck-ctrl
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
